# conf/supervisord.conf
[unix_http_server]
file = /var/log/jupo/supervisor.sock

[supervisord]
logfile       = /var/log/jupo/supervisord.log
directory     = /home/Workspace/jupo/src/
pidfile       = /var/log/jupo/supervisord.pid
nodaemon      = false 
minfds        = 65535                                    
minprocs      = 200                                 
environment   = PATH="/home/Workspace/jupo/env/bin" 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/log/jupo/supervisor.sock


[program:smtpd]
command       = python mail_server.py


[program:snowflake]
command       = python flake.py --worker_id=0 --port=2300


[program:push]
command       = python rqworker.py push
numprocs      = 8
process_name  = %(process_num)s
autorestart   = true


[program:fb_importer]
command       = python rqworker.py importer
numprocs      = 8
process_name  = %(process_num)s
autorestart   = true


[program:index]
command       = python rqworker.py index
numprocs      = 2
process_name  = %(process_num)s
autorestart   = true


[program:invite]
command       = python rqworker.py invite
numprocs      = 1
process_name  = %(process_num)s
autorestart   = true


[program:send_mail]
command       = python rqworker.py send_mail
numprocs      = 1
process_name  = %(process_num)s
autorestart   = true


[program:crawler]
command       = python rqworker.py urls
numprocs      = 1
process_name  = %(process_num)s
autorestart   = true


[program:move_to_s3]
command       = python rqworker.py move_to_s3
numprocs      = 1
process_name  = %(process_num)s
autorestart   = true


[program:notifier]
command       = python rqworker.py notification
numprocs      = 2
process_name  = %(process_num)s
autorestart   = true




[program:redis_server]
directory     = /home/Workspace/jupo/conf/
command       = redis-server redis_server.conf

[program:redis_pingpong_server]
directory     = /home/Workspace/jupo/conf/
command       = redis-server redis_pingpong_server.conf

[program:redis_pubsub_server]
directory     = /home/Workspace/jupo/conf/
command       = redis-server redis_pubsub_server.conf

[program:redis_taskqueue_server]
directory     = /home/Workspace/jupo/conf/
command       = redis-server redis_taskqueue_server.conf



[program:nginx]
command       = /usr/local/nginx/sbin/nginx -c /home/Workspace/jupo/conf/nginx.conf


[program:mongodb]
command       = mongod --dbpath=/data/mongo

[program:memcached]
command       = memcached -m 512 -p 11211 -l 0.0.0.0 -u memcached


[program:elasticsearch]
command       = /usr/local/elasticsearch/bin/elasticsearch -f


[program:pingpong_port7000]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:7000 main:app

[program:pingpong_port7001]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:7001 main:app

[program:pubsub_port8000]
command       = gunicorn -k gevent -w 2 --timeout 120 --max-requests 1000 -b 0.0.0.0:8000 main:app

[program:pubsub_port8001]
command       = gunicorn -k gevent -w 2 --timeout 120 --max-requests 1000 -b 0.0.0.0:8001 main:app

[program:pubsub_port8002]
command       = gunicorn -k gevent -w 2 --timeout 120 --max-requests 1000 -b 0.0.0.0:8002 main:app

[program:pubsub_port8003]
command       = gunicorn -k gevent -w 2 --timeout 120 --max-requests 1000 -b 0.0.0.0:8003 main:app

[program:web_port9000]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:9000 main:app

[program:web_port9001]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:9001 main:app

[program:web_port9002]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:9002 main:app

[program:web_port9003]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:9003 main:app

[program:web_port9500]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:9500 main:app

[program:web_port9501]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:9501 main:app

[program:web_port9502]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:9502 main:app

[program:web_port9503]
command       = gunicorn -k gevent -w 2 --max-requests 1000 -b 0.0.0.0:9503 main:app

